# REPOS/server/plex-media-server/plex
#!/bin/sh
### BEGIN INIT INFO
# Provides:          plexmediaserver
# Required-Start:    $remote_fs $networking
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Plex Media Server
# Description:       Plex Media Server for Linux,
#                    More information at http://www.plexapp.com
#                    Many thanks to the great PlexApp team for their wonderfull job !
# Author:            Cedric Quillevere / origin@killy.net
# Rewamped	     Christian Svedin / christian.svedin@gmail.com
# Version:           1.2
### END INIT INFO
. /lib/lsb/init-functions

BIN_FILE="/usr/lib/plexmediaserver/Plex Media Server"

#$LastChangedBy: dj $
#$Date: 2019-09-08 13:50:13 -0500 (Sun, 08 Sep 2019) $

case "$1" in
   start)
      log_info_msg "Starting plex..."
      start_daemon "${BIN_FILE}"&
      evaluate_retval
      ;;

   stop)
      	log_info_msg "Stopping plex..."
	killproc "${BIN_FILE}"
      	evaluate_retval
      ;;
   status)
       pidof "${BIN_FILE}"
       ;;
   restart)
      $0 stop
      sleep 1
      $0 start
      ;;
   *)
      echo "Usage: $0 {start|stop|restart|status}"
      exit 1
      ;;
esac

# End /etc/init.d/fcron
